# supervisord_coop_shaba_login.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Quotes around values are not supported, except in the case of
;    the environment= options as shown below.
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".
;  - Command will be truncated if it looks like a config file comment, e.g.
;    "command=bash -c 'foo ; bar'" will truncate to "command=bash -c 'foo ".
;
; Warning:
;  Paths throughout this example file use /tmp because it is available on most
;  systems.  You will likely need to change these to locations more appropriate
;  for your system.  Some systems periodically delete older files in /tmp.
;  Notably, if the socket file defined in the [unix_http_server] section below
;  is deleted, supervisorctl will be unable to connect to supervisord.

[unix_http_server]
file=./supervisor_coop_shaba_login.sock   ; the path to the socket file
chmod=0700                 ; socket file mode (default 0700)
username=admin              ; 设置用户名
password=123456            ; 设置密码

; Security Warning:
;  The inet HTTP server is not enabled by default.  The inet HTTP server is
;  enabled by uncommenting the [inet_http_server] section below.  The inet
;  HTTP server is intended for use within a trusted environment only.  It
;  should only be bound to localhost or only accessible from within an
;  isolated, trusted network.  The inet HTTP server does not support any
;  form of encryption.  The inet HTTP server does not use authentication
;  by default (see the username= and password= options to add authentication).
;  Never expose the inet HTTP server to the public internet.

[inet_http_server]         ; inet (TCP) server disabled by default
port=0.0.0.0:9003        ; ip_address:port specifier, *:port for all iface
username=admin              ; 设置用户名
password=123456            ; 设置密码

[supervisord]
logfile=./supervisord_coop_shaba_login.log ; main log file
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=./supervisord_coop_shaba_login.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
;umask=022                   ; process file creation umask; default 022
;user=supervisord            ; setuid to this UNIX account at startup; recommended if root
;identifier=supervisor       ; supervisord identifier, default is 'supervisor'
directory=%(here)s         ; 设置工作目录
;nocleanup=true              ; don't clean up tempfiles at start; default false
;childlogdir=/tmp            ; 'AUTO' child log dir, default $TEMP
;environment=KEY="value"     ; key value pairs to add to environment
;strip_ansi=false            ; strip ansi escape codes in logs; def. false

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.

[supervisorctl]
serverurl=unix://./supervisor_coop_shaba_login.sock ; use a unix:// URL  for a unix socket
serverurl=http://0.0.0.0:9003 ; use an http:// url to specify an inet socket
username=admin              ; 设置用户名
password=123456            ; 设置密码
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

; The sample program section below shows all possible program subsection values.
; Create one or more 'real' program: sections to be able to control them under
; supervisor.

;[program:theprogramname]
;command=/bin/cat              ; the program (relative uses PATH, can take args)
;process_name=%(program_name)s ; process_name expr (default %(program_name)s)
;numprocs=1                    ; number of processes copies to start (def 1)
;directory=/tmp                ; directory to cwd to before exec (def no cwd)
;umask=022                     ; umask for process (default None)
;priority=999                  ; the relative start priority (default 999)
;autostart=true                ; start at supervisord start (default: true)
;startsecs=1                   ; # of secs prog must stay up to be running (def. 1)
;startretries=3                ; max # of serial start failures when starting (default 3)
;autorestart=unexpected        ; when to restart if exited after running (def: unexpected)
;exitcodes=0                   ; 'expected' exit codes used with autorestart (default 0)
;stopsignal=QUIT               ; signal used to kill process (default TERM)
;stopwaitsecs=10               ; max num secs to wait b4 SIGKILL (default 10)
;stopasgroup=false             ; send stop signal to the UNIX process group (default false)
;killasgroup=false             ; SIGKILL the UNIX process group (def false)
;user=chrism                   ; setuid to this UNIX account to run the program
;redirect_stderr=true          ; redirect proc stderr to stdout (default false)
;stdout_logfile=/a/path        ; stdout log path, NONE for none; default AUTO
;stdout_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
;stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)
;stdout_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
;stdout_events_enabled=false   ; emit events on stdout writes (default false)
;stdout_syslog=false           ; send stdout to syslog with process name (default false)
;stderr_logfile=/a/path        ; stderr log path, NONE for none; default AUTO
;stderr_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
;stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)
;stderr_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
;stderr_events_enabled=false   ; emit events on stderr writes (default false)
;stderr_syslog=false           ; send stderr to syslog with process name (default false)
;environment=A="1",B="2"       ; process environment additions (def no adds)
;serverurl=AUTO                ; override serverurl computation (childutils)

; The sample eventlistener section below shows all possible eventlistener
; subsection values.  Create one or more 'real' eventlistener: sections to be
; able to handle event notifications sent by supervisord.

;[eventlistener:theeventlistenername]
;command=/bin/eventlistener    ; the program (relative uses PATH, can take args)
;process_name=%(program_name)s ; process_name expr (default %(program_name)s)
;numprocs=1                    ; number of processes copies to start (def 1)
;events=EVENT                  ; event notif. types to subscribe to (req'd)
;buffer_size=10                ; event buffer queue size (default 10)
;directory=/tmp                ; directory to cwd to before exec (def no cwd)
;umask=022                     ; umask for process (default None)
;priority=-1                   ; the relative start priority (default -1)
;autostart=true                ; start at supervisord start (default: true)
;startsecs=1                   ; # of secs prog must stay up to be running (def. 1)
;startretries=3                ; max # of serial start failures when starting (default 3)
;autorestart=unexpected        ; autorestart if exited after running (def: unexpected)
;exitcodes=0                   ; 'expected' exit codes used with autorestart (default 0)
;stopsignal=QUIT               ; signal used to kill process (default TERM)
;stopwaitsecs=10               ; max num secs to wait b4 SIGKILL (default 10)
;stopasgroup=false             ; send stop signal to the UNIX process group (default false)
;killasgroup=false             ; SIGKILL the UNIX process group (def false)
;user=chrism                   ; setuid to this UNIX account to run the program
;redirect_stderr=false         ; redirect_stderr=true is not allowed for eventlisteners
;stdout_logfile=/a/path        ; stdout log path, NONE for none; default AUTO
;stdout_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
;stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)
;stdout_events_enabled=false   ; emit events on stdout writes (default false)
;stdout_syslog=false           ; send stdout to syslog with process name (default false)
;stderr_logfile=/a/path        ; stderr log path, NONE for none; default AUTO
;stderr_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
;stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)
;stderr_events_enabled=false   ; emit events on stderr writes (default false)
;stderr_syslog=false           ; send stderr to syslog with process name (default false)
;environment=A="1",B="2"       ; process environment additions
;serverurl=AUTO                ; override serverurl computation (childutils)

; The sample group section below shows all possible group values.  Create one
; or more 'real' group: sections to create "heterogeneous" process groups.

;[group:thegroupname]
;programs=progname1,progname2  ; each refers to 'x' in [program:x] definitions
;priority=999                  ; the relative start priority (default 999)

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /opt/homebrew/etc/supervisor.d/*.ini

; 这是一个示例程序配置
# [program:example]
# command=python -m http.server 8000              ; 运行一个简单的 HTTP 服务器作为示例
# directory=%(here)s                             ; 程序的工作目录
# autostart=true                                ; 当 supervisord 启动时自动启动
# autorestart=true                              ; 当程序退出时自动重启
# redirect_stderr=true                          ; 将标准错误重定向到标准输出
# stdout_logfile=./example.log                  ; 日志文件位置
# environment=PYTHONUNBUFFERED="1"              ; 环境变量设置

; 定义监控工具组
# [group:monitors]
# programs=cpu_monitor,memory_monitor
# priority=999

# ; CPU监控脚本
# [program:cpu_monitor]
# command=python scripts/cpu_monitor.py
# directory=%(here)s
# autostart=true
# autorestart=true
# redirect_stderr=true
# stdout_logfile=./logs/cpu_monitor.log
# environment=PYTHONUNBUFFERED="1"

# ; 内存监控脚本
# [program:memory_monitor]
# command=python scripts/memory_monitor.py
# directory=%(here)s
# autostart=true
# autorestart=true
# redirect_stderr=true
# stdout_logfile=./logs/memory_monitor.log
# environment=PYTHONUNBUFFERED="1"

# ; 定义工具组
# [group:tools]
# programs=time_printer,counter
# priority=900

# ; 时间打印脚本
# [program:time_printer]
# command=python scripts/time_printer.py
# directory=%(here)s
# autostart=true
# autorestart=true
# redirect_stderr=true
# stdout_logfile=./logs/time_printer.log
# environment=PYTHONUNBUFFERED="1"

# ; 计数器脚本
# [program:counter]
# command=python scripts/counter.py
# directory=%(here)s
# autostart=true
# autorestart=true
# redirect_stderr=true
# stdout_logfile=./logs/counter.log
# environment=PYTHONUNBUFFERED="1"

[program:mybc_pinnacle_relogin_bet_check]
command=python prod_v3/pinnacle_relogin_bet_check.py
directory=/Users/feiwentao/tianhei_projects/python_projests/bc2
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/Users/feiwentao/tianhei_projects/python_projests/bc2/logs/pinnacle_relogin_bet_check.log
stderr_logfile=/Users/feiwentao/tianhei_projects/python_projests/bc2/logs/pinnacle_relogin_bet_check.err.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=10
startretries=3

[program:mybc_shaba_relogin_bet_check]
command=python prod_v3/shaba_relogin_bet_check.py
directory=/Users/feiwentao/tianhei_projects/python_projests/bc2
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/Users/feiwentao/tianhei_projects/python_projests/bc2/logs/shaba_relogin_bet_check.log
stderr_logfile=/Users/feiwentao/tianhei_projects/python_projests/bc2/logs/shaba_relogin_bet_check.err.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=10
startretries=3


[program:mybc_wanmei_relogin_check]
command=python prod_v3/wanmei_relogin_check.py
directory=/Users/feiwentao/tianhei_projects/python_projests/bc2
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/Users/feiwentao/tianhei_projects/python_projests/bc2/logs/wanmei_relogin_check.log
stderr_logfile=/Users/feiwentao/tianhei_projects/python_projests/bc2/logs/wanmei_relogin_check.err.log
environment=PYTHONUNBUFFERED=1
stopasgroup=true
killasgroup=true
startsecs=10
startretries=3


[group:mybc_login_check]
programs=mybc_pinnacle_relogin_bet_check,mybc_shaba_relogin_bet_check,mybc_wanmei_relogin_check




