# conf/subpro/mytesta.ini
[program:mytestproa]
command=/bin/python mytest.py              ; the program (relative uses PATH, can take args)
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
;numprocs=1                    ; number of processes copies to start (def 1)
directory=/home/                ; directory to cwd to before exec (def no cwd)
comments=我的测试程序a
