# supervisor/tests/fixtures/issue-291a.conf
[supervisord]
loglevel=debug               ; log level; default info; others: debug,warn,trace
logfile=/tmp/issue-291a.log  ; main log file; default $CWD/supervisord.log
pidfile=/tmp/issue-291a.pid  ; supervisord pidfile; default supervisord.pid
nodaemon=true                ; start in foreground if true; default false

[program:print_env]
command=python %(here)s/print_env.py
startsecs=0
autorestart=false
