# supervisor/tests/fixtures/issue-1224.conf
[supervisord]
nodaemon = true
pidfile = /tmp/issue-1224.pid
nodaemon = true
logfile = /dev/stdout
logfile_maxbytes = 0

[program:cat]
command = /bin/cat
startsecs = 0

