# supervisor/tests/fixtures/issue-1298.conf
[supervisord]
nodaemon=true                ; start in foreground if true; default false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/issue-1298.sock    ; the path to the socket file

[supervisorctl]
serverurl=unix:///tmp/issue-1298.sock   ; use a unix:// URL  for a unix socket

[program:spew]
command=python %(here)s/spew.py
