# supervisor/tests/fixtures/issue-1170a.conf
[supervisord]
nodaemon=true                ; start in foreground if true; default false
loglevel=debug                ; log level; default info; others: debug,warn,trace
logfile=/tmp/issue-1170a.log  ; main log file; default $CWD/supervisord.log
pidfile=/tmp/issue-1170a.pid  ; supervisord pidfile; default supervisord.pid
environment=FOO="set from [supervisord] section"

[program:echo]
command=bash -c "echo '%(ENV_FOO)s'"
startsecs=0
startretries=0
autorestart=false
