# supervisor/tests/fixtures/issue-1184.conf
[supervisord]
loglevel = debug
logfile=/tmp/issue-1184.log
pidfile=/tmp/issue-1184.pid
nodaemon = true

[unix_http_server]
file = /tmp/issue-1184.sock

[supervisorctl]
serverurl = unix:///tmp/issue-1184.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:cat]
command = /bin/cat
startsecs = 0

