# supervisor/tests/fixtures/issue-986.conf
[supervisord]
loglevel=debug               ; log level; default info; others: debug,warn,trace
logfile=/tmp/issue-986.log   ; main log file; default $CWD/supervisord.log
pidfile=/tmp/issue-986.pid   ; supervisord pidfile; default supervisord.pid
nodaemon=true                ; start in foreground if true; default false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/issue-986.sock   ; the path to the socket file

[supervisorctl]
serverurl=unix:///tmp/issue-986.sock  ; use a unix:// URL  for a unix socket

[program:echo]
command=bash -c "echo 'dhcrelay -d -q -a %%h:%%p %%P -i Vlan1000 192.168.0.1'"
startsecs=0
autorestart=false
