# supervisor/tests/fixtures/issue-565.conf
[supervisord]
loglevel=info                ; log level; default info; others: debug,warn,trace
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true                ; start in foreground if true; default false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:hello]
command=bash supervisor/tests/fixtures/hello.sh
stdout_events_enabled=true

[eventlistener:listener]
command=python supervisor/tests/fixtures/listener.py
events=PROCESS_LOG
