# supervisor/tests/fixtures/issue-1054.conf
[supervisord]
loglevel = debug
logfile=/tmp/issue-1054.log
pidfile=/tmp/issue-1054.pid
nodaemon = true

[unix_http_server]
file=/tmp/issue-1054.sock     ; the path to the socket file

[supervisorctl]
serverurl=unix:///tmp/issue-1054.sock ; use a unix:// URL  for a unix socket

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:cat]
command = /bin/cat
startsecs = 0

