# supervisor/tests/fixtures/issue-664.conf
[supervisord]
loglevel=debug
logfile=/tmp/issue-664.log
pidfile=/tmp/issue-664.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/issue-664.sock     ; the path to the socket file

[supervisorctl]
serverurl=unix:///tmp/issue-664.sock ; use a unix:// URL  for a unix socket

[program:test_öäü]
command = /bin/cat
startretries = 0
autorestart = false

