# supervisor/tests/fixtures/issue-836.conf
[supervisord]
loglevel = debug
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
nodaemon = true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/issue-565.sock     ; the path to the socket file

[supervisorctl]
serverurl=unix:///tmp/issue-565.sock ; use a unix:// URL  for a unix socket

[program:cat]
command = /bin/cat
startretries = 0
autorestart = false

