# supervisor/tests/fixtures/issue-550.conf
[supervisord]
loglevel=info                ; log level; default info; others: debug,warn,trace
logfile=/tmp/issue-550.log   ; main log file; default $CWD/supervisord.log
pidfile=/tmp/issue-550.pid   ; supervisord pidfile; default supervisord.pid
nodaemon=true                ; start in foreground if true; default false
environment=THIS_SHOULD=BE_IN_CHILD_ENV

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/issue-550.sock     ; the path to the socket file

[supervisorctl]
serverurl=unix:///tmp/issue-550.sock ; use a unix:// URL  for a unix socket

[program:print_env]
command=python %(here)s/print_env.py
startsecs=0
startretries=0
autorestart=false
