# supervisor/tests/fixtures/issue-638.conf
[supervisord]
loglevel=debug               ; log level; default info; others: debug,warn,trace
logfile=/tmp/issue-638.log   ; main log file; default $CWD/supervisord.log
pidfile=/tmp/issue-638.pid   ; supervisord pidfile; default supervisord.pid
nodaemon=true                ; start in foreground if true; default false

[program:produce-unicode-error]
command=bash -c 'echo -e "\x88"'
startretries=0
autorestart=false

